# templates/calico-node.service.erb
<%# Inspired from https://github.com/projectcalico/calico-containers/blob/master/docs/CalicoAsService.md -%>
[Unit]
Description=Calico node service
After=docker.service etcd2.service
Requires=docker.service etcd2.service

[Service]
ExecStartPre=-/usr/bin/docker rm -f calico-node
ExecStart=/opt/bin/calicoctl node run --init-system --ip6=<%= @ipaddress6 %>
ExecStop=-/usr/bin/docker stop calico-node

[Install]
WantedBy=multi-user.target
