# templates/puppet.service.erb
[Unit]
Description=Puppet agent in Docker
After=docker.service
Requires=docker.service

[Service]
ExecStartPre=/bin/bash -c '/usr/bin/docker rm -f %n || true'
ExecStop=/bin/bash -c '/usr/bin/docker rm -f %n || true'
<% if @has_ipmi %>
ExecStartPre=/bin/bash -c 'set -e -x; modprobe ipmi_si; rmmod ipmi_devintf || true; rm -rf /dev/ipmi0; modprobe ipmi_devintf; for i in $(seq 1 100); do test -f /dev/ipmi0 && break; sleep 0.1; done'
<% end %>
ExecStart=/usr/bin/docker run <% -%>
  --name %n <% -%>
  --net=host <% -%>
  --privileged <% -%>
  -v /:/opt/root <% -%>
  -v /etc/systemd:/etc/systemd <% -%>
  -v /etc/puppet:/etc/puppet <% -%>
  -v /var/lib/puppet:/var/lib/puppet <% -%>
  -v /home/core:/home/core <% -%>
  -v /etc/os-release:/etc/os-release:ro <% -%>
  -v /etc/lsb-release:/etc/lsb-release:ro <% -%>
  -v /etc/coreos:/etc/coreos:rw <% -%>
  -v /run:/run:ro <% -%>
  -v /usr/bin/systemctl:/usr/bin/systemctl:ro <% -%>
  -v /lib64:/lib64:ro <% -%>
  -v /usr/lib64/systemd:/usr/lib64/systemd <% -%>
  -v /usr/lib/systemd:/usr/lib/systemd <% -%>
  -v /sys/fs/cgroup:/sys/fs/cgroup:ro <% -%>
<% if @has_ipmi -%>
  -v /dev/ipmi0:/dev/ipmi0 <% -%>
<% end -%>
<% @facts.each do |name, value| -%>
  -e "FACTER_<%= name %>=<%= value %>" <% -%>
<% end -%>
  epflsti/cluster.coreos.puppet:latest <% -%>
  agent --no-daemonize --logdest=console --environment=production
RestartSec=5s
Restart=always

[Install]
WantedBy=multi-user.target
