# templates/puppet.service.erb
[Unit]
Description=Puppet agent in Docker
After=docker.service
Requires=docker.service

[Service]
ExecStartPre=/bin/bash -c '/usr/bin/docker rm -f %n || true'
ExecStartPre=/bin/rm -f /var/lib/puppet/state/agent_catalog_run.lock
ExecStop=/bin/bash -c '/usr/bin/docker rm -f %n || true'
ExecStart=/usr/bin/docker run <% -%>
  --name %n <% -%>
  --net=host <% -%>
  --privileged <% -%>
  -v /:/opt/root <% -%>
  -v /dev:/dev <% -%>
  -v /etc/systemd:/etc/systemd <% -%>
  -v /etc/ssh:/etc/ssh <% -%>
  -v /etc/puppet:/etc/puppet <% -%>
  -v /var/lib/puppet:/var/lib/puppet <% -%>
  -v /home/core:/home/core <% -%>
  -v /etc/os-release:/etc/os-release:ro <% -%>
  -v /etc/lsb-release:/etc/lsb-release:ro <% -%>
  -v /etc/coreos:/etc/coreos:rw <% -%>
  -v /run:/run:ro <% -%>
  -v /usr/bin/systemctl:/usr/bin/systemctl:ro <% -%>
  -v /lib64:/lib64:ro <% -%>
  -v /usr/lib64/systemd:/usr/lib64/systemd <% -%>
  -v /usr/lib/systemd:/usr/lib/systemd <% -%>
  -v /sys/fs/cgroup:/sys/fs/cgroup:ro <% -%>
<% @facts.each do |name, value| -%>
  -e "FACTER_<%= name %>=<%= value %>" <% -%>
<% end -%>
  epflsti/cluster.coreos.puppet:latest <% -%>
  agent --no-daemonize --logdest=console --environment=production
RestartSec=5s
Restart=always

[Install]
WantedBy=multi-user.target
