# templates/vipv4.service.erb
[Unit]
Description=ucarp on <%= @where == "internal" ? "external" : "gateway" %> VIPv4 <%= @ip %>
Requires=docker.service network.target
After=docker.service network.target

[Service]
Restart=always
RestartSec=60s
TimeoutStartSec=120
TimeoutStopSec=25

# See https://github.com/nicolerenee/docker-ucarp/blob/master/start.sh
# for how the environment variables ought to be set
# UCARP_VHID and UCARP_VIRTUALADDRESS are defined as Environment= stanzas
# in a supplemental, per-service systemd config file e.g.
# /etc/systemd/system/<%= @cluster_owner %>.gateway-ucarp-external@MYVIP.service.d/ip-and-vhid.conf
ExecStart=/usr/bin/docker run \
    --rm \
    --name %p  \
    --net=host --cap-add NET_ADMIN \
    nicolerenee/ucarp \
      /usr/sbin/ucarp \
<%
# From the ยง in ucarp/README that starts with "To understand how ucarp
# works," we see that we need to tweak --advbase to get a "sticky"
# setting for the preferred master:

if @vipv4_affinity_table && @vipv4_affinity_table[@ip]; then
   @index = @vipv4_affinity_table[@ip].index(@hostname)
   @priority = @index ? @index + 1 : 30
else
   @priority = 30
end

# 
@all_vips = (
      scope['::epflsti_coreos::private::gateway::ucarp::external_ipv4_vips'] + 
      scope['::epflsti_coreos::private::gateway::ucarp::gateway_ipv4_vips'] )
-%>
      --advbase=<%= @priority %> \
<%- if @priority == 1 -%>
      --preempt \
<%- end -%>
      --interface=<%= @interface %> \
      --srcip=<%= @membership_protocol_ip %> \
      --vhid=<%= 100 + @all_vips.index(@ip) %> \
      --pass=<%= @failover_shared_secret %> \
      --addr=<%= @ip %> \
      --upscript=/etc/ucarp/vip-up-default.sh \
      --downscript=/etc/ucarp/vip-down-default.sh \
      --nomcast

ExecStop=/usr/bin/docker stop %p

